# scripts/systemd/sausalito-backup.timer
[Unit]
Description=Run Sausalito postgres backup guard daily

[Timer]
OnCalendar=*-*-* 04:10:00
Persistent=true

[Install]
WantedBy=timers.target
