# scripts/systemd/sausalito-certbot-renew.timer
[Unit]
Description=Run Sausalito certbot renew daily

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
