# scripts/systemd/sausalito-disk-guard.timer
[Unit]
Description=Run Sausalito disk guard every 6 hours

[Timer]
OnCalendar=*-*-* 00,06,12,18:20:00
Persistent=true

[Install]
WantedBy=timers.target
